# deploy/tasks/crontab

### variables ##
  SHELL=/bin/bash
  PATH=:/bin:/usr/bin:/usr/local/bin:/usr/sbin:/usr/local/sbin
  MAILTO=john.skilbeck@gmail.com
  cmd="deploy/bin/run-job -m jobs.stock-market"
  app_dir="/home/skilbjo/deploy/app/router-logs"

## jobs ##
5 * * * * skilbjo . /home/skilbjo/.profilelocal; cd "$app_dir" ; $cmd >/dev/null
